# roles/docker_service/templates/batch.timer.j2
[Unit]
Description=Run {{ item.name }}.service on a schedule

[Timer]
OnCalendar={{ item.schedule }}

[Install]
WantedBy=multi-user.target
